# factorytest_script/99-usb-config.rules
SUBSYSTEMS=="usb", ACTION=="add", KERNEL=="sd*[!0-9]", RUN+="/usr/bin/usb_config.sh %k"